# fpm/fpm_source/usr/local/dumbo/etc/apt/journald-cloudwatch-logs.postinst
if systemctl show-environment &>/dev/null; then
    systemctl daemon-reload
    systemctl enable journald-cloudwatch-logs.service
    systemctl reload-or-restart journald-cloudwatch-logs
fi
